# docs/static/vmware-reporter.service
[Unit]
Description=Run vmware-reporter (all)

[Service]
Type=oneshot
WorkingDirectory=/srv/vmware-reporter
ExecStart=vmware-reporter all
